# dockerfiles/s6-rc.d/ssh/run
#!/command/execlineb -P
with-contenv
s6-setuidgid plural
plural serve