# dockerfiles/s6-rc.d/dockerd/run
#!/command/with-contenv sh
exec dockerd --log-level error